# root.cron
0,30 * * * *	FORCE_IMAGE_REMOVAL=1 /docker-gc
